# payload/Library/LaunchAgents/com.github.go-installapplications.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>com.github.go-installapplications.agent</string>

        <key>Program</key>
        <string>/Library/go-installapplications/go-installapplications</string>

        <key>ProgramArguments</key>
        <array>
            <string>/Library/go-installapplications/go-installapplications</string>
            <string>--mode</string>
            <string>agent</string>
        </array>

        <key>RunAtLoad</key>
        <true/>

        <key>StandardOutPath</key>
        <string>/var/log/go-installapplications/go-installapplications.agent.log</string>

        <key>StandardErrorPath</key>
        <string>/var/log/go-installapplications/go-installapplications.agent.log</string>
    </dict>
</plist>